# deploy/systemd/molodost48-api.service
[Unit]
Description=molodost48 API (Fastify + Prisma)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=molodost
Group=molodost
WorkingDirectory=/var/www/molodost48/api
EnvironmentFile=/etc/molodost48/api.env
ExecStart=/usr/bin/node --import tsx/esm src/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=molodost48-api

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/molodost48/uploads /var/lib/molodost48/backups /var/log/molodost48
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits (1GB total RAM)
MemoryMax=384M
TasksMax=64

[Install]
WantedBy=multi-user.target